# install/etc/init/addacs-pinmux.conf
# addacs-pinmux
#
# 

description	"ADDACS pin mux for beagleboard"

start on starting addacs



#expect daemon
#respawn

script
	
end script

