# install/etc/init/addacs-dbus.conf

description "create a dbus for addacs"

start on starting addacs

script
	dbus-launch --auto-syntax > /tmp/dbusbits
	echo "wrote file " >> /tmp/thing.log
end script
